# extra/files/etc/uci-defaults/80_cfg_enable_openssh_root
#!/bin/sh
#
# Copyright (C) 2016 honwen.chan@gmail.com
#

[ -f /etc/ssh/sshd_config ] || exit0

sed 's/.*PermitRootLogin.*/PermitRootLogin yes/g' -i /etc/ssh/sshd_config

exit 0
