# extra/files/etc/uci-defaults/90_cfg_enable_miniupnp
#!/bin/sh
#
# Copyright (C) 2017 honwen.chan@gmail.com
#

[ -f /etc/config/upnpd ] || exit0

uci batch <<EOF
	set upnpd.config.enabled=1
	commit upnpd
EOF
